# dot_config/systemd/user/logrotate-nvim.timer
[Unit]
Description=Daily logrotate for Neovim LSP

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
